# tapo-temperature.service
[Unit]
Description=Система контролю температури для твердопаливного котла
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/tapo
Environment="PATH=/home/pi/tapo/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/home/pi/tapo/venv/bin/python /home/pi/tapo/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Безпека
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

