# temperature.service
[Unit]
Description=Система контролю температури для твердопаливного котла (Sonoff S60)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Замініть 'raccoon' на ваше ім'я користувача або 'pi' для Raspberry Pi OS
User=raccoon
Group=raccoon
WorkingDirectory=/home/raccoon/rpi-temp-controll
Environment="PATH=/home/raccoon/rpi-temp-controll/venv/bin"
ExecStart=/home/raccoon/rpi-temp-controll/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Дозволити доступ до GPIO та SPI
SupplementaryGroups=gpio spi

# Безпека
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

